# vnc-monitoring/vnc-monitor.service
[Unit]
Description=VNC Server Monitoring Agent
Documentation=man:vnc_monitor.py
After=network.target
Wants=network.target

[Service]
Type=simple
User=root
Group=root
WorkingDirectory=/opt/vnc-monitoring
ExecStart=/usr/bin/python3 /opt/vnc-monitoring/vnc_monitor.py --daemon --display :1 --interval 30 --log-file /var/log/vnc_monitor.log
ExecReload=/bin/kill -HUP $MAINPID
ExecStop=/bin/kill -TERM $MAINPID
Restart=always
RestartSec=10
StandardOutput=journal
StandardError=journal
SyslogIdentifier=vnc-monitor

# Security settings
NoNewPrivileges=true
PrivateTmp=true
ProtectSystem=strict
ProtectHome=true
ReadWritePaths=/var/log /opt/vnc-monitoring
CapabilityBoundingSet=CAP_KILL CAP_SETUID CAP_SETGID CAP_NET_BIND_SERVICE

# Resource limits
LimitNOFILE=1024
LimitNPROC=64

[Install]
WantedBy=multi-user.target